# ubuntu-noble-lts-amd64-openstack-nova/ansible/roles/nova/templates/etc/systemd/system/nova-migrate.timer
# {{ ansible_managed }}

[Unit]
Description="nova-migrate"

[Timer]
OnBootSec=5min
OnUnitActiveSec=24h
Unit=nova-migrate.service

[Install]
WantedBy=multi-user.target
