# ubuntu-noble-lts-amd64-openstack-cinder/ansible/roles/cinder/templates/etc/systemd/system/cinder-migrate.timer
# {{ ansible_managed }}

[Unit]
Description="cinder-migrate"

[Timer]
OnBootSec=5min
OnUnitActiveSec=24h
Unit=cinder-migrate.service

[Install]
WantedBy=multi-user.target
