# ubuntu-noble-lts-amd64-openstack-cinder/ansible/roles/cinder/templates/etc/systemd/system/apache2.service.d/override.conf
# {{ ansible_managed }}

[Unit]

# Need to make sure cinder Step Certs are generated.
Wants=cinder-step-certs.service

# cinder Step Certs need to happen first.
After=cinder-step-certs.service
