# ubuntu-noble-lts-amd64-openstack-cinder/ansible/roles/cinder/templates/etc/systemd/system/cinder-scheduler.service.d/override.conf
# {{ ansible_managed }}

[Unit]

Wants=cinder-migrate.service
After=cinder-migrate.service
