# ubuntu-noble-lts-amd64-openstack-nova/ansible/roles/nova/templates/etc/systemd/system/nova-scheduler.service.d/override.conf
# {{ ansible_managed }}

[Unit]
Wants=nova-migrate.service
After=nova-migrate.service
