# ubuntu-noble-lts-amd64-patroni/ansible/roles/patroni/templates/etc/systemd/system/patroni.service.d/override.conf
# {{ ansible_managed }}

[Unit]

# Need to make sure Patroni Step Certs are generated.
Wants=patroni-step-certs.service

# Patroni Step Certs need to happen first.
After=patroni-step-certs.service
