# ubuntu-noble-lts-amd64-openstack-octavia/ansible/roles/octavia/templates/etc/systemd/system/octavia-api.service.d/override.conf
# {{ ansible_managed }}

[Unit]
Wants=haproxy.service octavia-migrate.service
After=haproxy.service octavia-migrate.service