# ubuntu-noble-lts-amd64-openstack-glance/ansible/roles/glance/templates/etc/systemd/system/glance-migrate.timer
# {{ ansible_managed }}

[Unit]
Description="glance-migrate"

[Timer]
OnBootSec=5min
OnUnitActiveSec=24h
Unit=glance-migrate.service

[Install]
WantedBy=multi-user.target
