# ubuntu-noble-lts-amd64-devstack/ansible/roles/devstack/templates/etc/systemd/system/tailscale-up.timer
# {{ ansible_managed }}

[Unit]
Description="tailscale-up"

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=tailscale-up.service

[Install]
WantedBy=multi-user.target
