# ubuntu-noble-lts-amd64-openstack-nova/ansible/roles/nova/templates/etc/systemd/system/nova-api.service.d/override.conf
# {{ ansible_managed }}

[Unit]
Wants=haproxy.service nova-migrate.service
After=haproxy.service nova-migrate.service
