# ubuntu-noble-lts-amd64-openstack-glance/ansible/roles/glance/templates/etc/systemd/system/glance-api.service.d/override.conf
# {{ ansible_managed }}

[Unit]
Wants=haproxy.service glance-migrate.service
After=haproxy.service glance-migrate.service
