# ubuntu-noble-lts-amd64-openstack-neutron/ansible/roles/neutron/templates/etc/systemd/system/haproxy.service.d/override.conf
# {{ ansible_managed }}

[Unit]

# Need to make sure HAProxy Step Certs are generated.
Wants=haproxy-step-certs.service

# HAProxy Step Certs need to happen first.
After=haproxy-step-certs.service
