# ubuntu-noble-lts-amd64-openstack-keystone/ansible/roles/keystone/templates/etc/systemd/system/keystone-project-users-token.service
# {{ ansible_managed }}

[Unit]
Description="keystone-project-users-token"
Requires=network-online.target

After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=no

User=keystone
Group=keystone

ProtectSystem=full
ProtectHome=read-only
PrivateTmp=yes

ExecStart=/usr/local/bin/keystone-project-users-token

[Install]
WantedBy=multi-user.target
