# ubuntu-noble-lts-amd64-openstack-keystone/ansible/roles/keystone/templates/etc/systemd/system/keystone-service-users-token.timer
# {{ ansible_managed }}

[Unit]
Description="keystone-service-users-token"

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=keystone-service-users-token.service

[Install]
WantedBy=multi-user.target
