# ubuntu-noble-lts-amd64-openstack-keystone/ansible/roles/keystone/templates/etc/systemd/system/keystone-endpoint-register.timer
# {{ ansible_managed }}

[Unit]
Description="keystone-endpoint-register"

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=keystone-endpoint-register.service

[Install]
WantedBy=multi-user.target
