# ubuntu-noble-lts-amd64-openstack-neutron/ansible/roles/neutron/templates/etc/systemd/system/neutron-server.service.d/override.conf
# {{ ansible_managed }}

[Unit]
Wants=haproxy.service neutron-migrate.service
After=haproxy.service neutron-migrate.service
