# ubuntu-noble-lts-amd64-openstack-placement/ansible/roles/placement/templates/etc/systemd/system/placement-migrate.timer
# {{ ansible_managed }}

[Unit]
Description="placement-migrate"

[Timer]
OnBootSec=5min
OnUnitActiveSec=24h
Unit=placement-migrate.service

[Install]
WantedBy=multi-user.target
