# ubuntu-noble-lts-amd64-openstack-octavia/ansible/roles/octavia/templates/etc/systemd/system/octavia-health-manager.service.d/override.conf
# {{ ansible_managed }}

[Unit]
Wants=octavia-migrate.service
After=octavia-migrate.service